# etc/rc.d/init.d/aqd.rh
#!/bin/bash
#
# aqd	control an Aquilon broker
#
### BEGIN INIT INFO
# Provides: aqd
# Required-Start: $local_fs $remote_fs
# Required-Stop: $local_fs $remote_fs
# Should-Start: $syslog $network
# Should-Stop: $syslog $network
# Default-Start: 2 3 5
# Default-Stop: 0 1 6
# Short-Description: Aquilon Broker
# Description: Quattor Aquilon Broker Daemon
### END INIT INFO

# Source function library
. /etc/rc.d/init.d/functions
. /etc/sysconfig/aquilon

retval=0

start() {
    if [[ -x $PY_BIN ]]; then
        echo -n "Starting aquilon broker: "
        if [[ $UID -ne 0 ]]; then
            retval=1
            failure
        else
            $PY_BIN $TWISTD --logfile=$LOGFILE --pidfile=$PIDFILE aqd --config=$conffile
            retval=$?
            if [[ $retval -eq 0 ]]; then
                touch $lockfile
                success
            fi
        fi
        echo
        return $retval
    else
        exit 5
    fi
}

stop() {
    echo -n "Stopping aquilon broker: "
    if [ $UID -ne 0 ]; then
        retval=1
        failure
    else
        killproc -p $PIDFILE $PY_BIN
        retval=$?
        if [[ $retval -eq 0 ]]; then
            rm -f $lockfile
        fi
    fi
    echo
    return $retval
}

restart(){
    stop
    start
}

aqstatus(){
    status aqd >/dev/null 2>&1
}

# See how we were called.
case "$1" in
    start)
        aqstatus && exit 0
        start
        retval=$?
        ;;
    stop)
        aqstatus || exit 0
        stop
        retval=$?
        ;;
    restart)
        restart
        retval=$?
        ;;
    force-reload)
        restart
        retval=$?
        ;;
    status)
        status -p $PIDFILE aqd
        retval=$?
        ;;
    *)
        echo $"Usage: $0 {start|stop|status|restart|force-reload}"
        retval=2
esac

exit $retval
